# ddns.service
[Unit]
Description=Update DNS

[Service]
ExecStart=/usr/bin/dnspodc --login_token $TOKEN --domain $DOMAIN --subdomain $SUBDOMAIN


